# SupervisorConf/sample.ini
[program:kubon]
directory=""
command=""
autostart=True
autorestart=False
stderr_logfile=/tmp/kubon/kubon_stderr.log
stdout_logfile=/tmp/kubon/kubon_stdout.log
user=root